# rc/.config/systemd/user/serve.service
[Unit]
Description=Media

[Service]
ExecStart=%h/code/serve/bin/serve %h/media

[Install]
WantedBy=default.target
